# k.conf
[program:k]
command=python3 /home/curator2/k/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/k/test.err.log
stdout_logfile=/home/curator2/k/test.out.log